# crontab
*/5 * * * * cd /app/ && python3 -m scraper.main > /proc/1/fd/1 2> /proc/1/fd/2
